# ports-overlay/trident/trident-core/pkg-plist
/boot/loader.conf.local
/etc/rc.conf.d/trident.conf
/etc/sysctl.d/trident-core.conf
bin/check-tor-mode
bin/disable-tor-mode
bin/enable-tor-mode
bin/mk_zfs_homedir
etc/asound.conf
etc/devd/trident-media.conf
etc/fonts/local.conf.trident
etc/init.d/net-online
etc/init.d/trident_init
etc/luminaDesktop.conf
etc/openal/alsoft.conf
etc/pcdm.conf.trident
etc/sudoers.d/trident
etc/trident/ipfw.tor.rules
etc/trident/recommended-services
etc/trident/required-services
etc/trident/torrc
share/cups/data/default-testpage.orig.pdf
share/fluxbox/styles/bora_blue_4k/theme.cfg
share/icons/default/index.theme
share/icons/hicolor/256x256/places/start-here-lumina.png
share/trident/pcdm-theme/trident.theme
share/trident/scripts/detect-best-driver.sh
share/trident/scripts/generate-desktop-entry
share/trident/scripts/generate-xorg-conf.sh
share/trident/scripts/remove-desktop-entry
share/trident/scripts/reset-firewall
share/trident/scripts/reset-firewall-tor
share/trident/scripts/sys-init.sh
share/trident/scripts/update-xorg-fonts.sh
share/trident/scripts/validate-services.sh
share/trident/xorg-templates/xorg.conf
share/trident/xstartup-parallel/disable-beep.sh
share/trident/xstartup-parallel/load-xresources.sh
share/trident/xstartup/01_xdg-user-dirs.sh
share/trident/xstartup/02_fix-gtk-environ.sh
share/trident/xstartup/enable-inputmethod.sh
share/trident/xstartup/gpg-agent.sh
share/trident/xstartup/ssh-agent.sh
share/wallpapers/trident/Trident-wallpaper-corner-darkgrey.png
share/wallpapers/trident/Trident-wallpaper-corner-lightgrey.png
share/wallpapers/trident/Trident-wallpaper-darkgrey.png
share/wallpapers/trident/Trident-wallpaper-lightgrey.png
share/wallpapers/trident/Trident-wallpaper2-corner-darkgrey.png
share/wallpapers/trident/Trident-wallpaper2-corner-lightgrey.png
share/wallpapers/trident/Trident-wallpaper2-darkgrey.png
share/wallpapers/trident/Trident-wallpaper2-lightgrey.png
share/wallpapers/trident/Trident-wallpaper3-corner-darkgrey.png
share/wallpapers/trident/Trident-wallpaper3-corner-lightgrey.png
share/wallpapers/trident/Trident-wallpaper3-darkgrey.png
share/wallpapers/trident/Trident-wallpaper3-lightgrey.png
share/wallpapers/trident/Trident-wallpaper4-corner-darkgrey.png
share/wallpapers/trident/Trident-wallpaper4-corner-lightgrey.png
share/wallpapers/trident/Trident-wallpaper4-darkgrey.png
share/wallpapers/trident/Trident-wallpaper4-lightgrey.png
share/wallpapers/trident/default.png
